# launchd/com.melvin.fresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.melvin.fresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/melvin/Code/dotfiles/fresh.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/melvin/Code/dotfiles</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>ThrottleInterval</key>
    <integer>259200</integer>

    <key>StandardOutPath</key>
    <string>/Users/melvin/Library/Logs/fresh-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/melvin/Library/Logs/fresh-launchd.error.log</string>
</dict>
</plist>
